# systemd/volatile-tmp.timer
[Unit]
Description=Run cleanup volatile tmp service

[Timer]
OnCalendar=daily
Persistent=true     
 
[Install]
WantedBy=timers.target